# deploy/systemd/vocabbook-cloud.service.example
[Unit]
Description=VocabBook Cloud API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/vocabbook-cloud/cloud_server
Environment=DATABASE_URL=sqlite+aiosqlite:////var/www/vocabbook-cloud/cloud_server/cloud_app.db
Environment=SECRET_KEY=CHANGE_ME_TO_A_LONG_RANDOM_SECRET
Environment=ALIPAY_APP_ID=CHANGE_ME
Environment=ALIPAY_PRIVATE_KEY_PATH=/var/www/vocabbook-cloud/secure/alipay_private_key.pem
Environment=ALIPAY_PUBLIC_KEY_PATH=/var/www/vocabbook-cloud/secure/alipay_public_key.pem
Environment=ALIPAY_NOTIFY_URL=https://api.historyai.fun/api/pay/alipay/notify
ExecStart=/var/www/vocabbook-cloud/.venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8010
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
